# debian/postinst
#!/bin/sh -e


case "$1" in
    configure|upgrade)
        WEBSERVERS="apache apache-ssl apache2"
        for apache in $WEBSERVERS; do
            if [ -d "/etc/$apache/conf.d" ] && \
               [ -x "/etc/init.d/$apache" ] && \
               [ ! -e "/etc/$apache/conf.d/p4a.conf" ]
            then
                echo "Installing p4a configuration file for $apache.."
                ln -sf /etc/p4a/apache.conf \
                    /etc/$apache/conf.d/p4a.conf
                if [ -f "/var/run/$apache.pid" ]; then
                    invoke-rc.d $apache reload
                fi
            fi
        done
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac


# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

